# launchd/com.local.wispr-freeze.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.local.wispr-freeze</string>
    <key>ProgramArguments</key>
    <array>
        <string>$HOME/Library/Scripts/WisprThrottle/check-wispr-freeze.sh</string>
    </array>
    <key>StartInterval</key>
    <integer>15</integer>
    <key>EnvironmentVariables</key>
    <dict>
        <key>WISPR_TIMEOUT</key>
        <string>5</string>
    </dict>
    <key>RunAtLoad</key>
    <true/>
</dict>
</plist>
